# vfd.service
[Unit]
Description=VFD daemon

[Service]
PIDFile=vfd.pid
ExecStart=/usr/bin/vfd_display.sh

[Install]
WantedBy=multi-user.target
